# Firmware/etc/udev/rules.d/60-persistent-v4l.rules
# do not edit this file, it will be overwritten on update

ACTION=="remove", GOTO="persistent_v4l_remove"
SUBSYSTEM!="video4linux", GOTO="persistent_v4l_end"
ENV{MAJOR}=="", GOTO="persistent_v4l_end"

IMPORT{program}="v4l_id $devnode"

SUBSYSTEMS=="usb", IMPORT{builtin}="usb_id"
IMPORT{builtin}="path_id"

# USB HUB: DM1/DP1
#ENV{ID_PATH}=="platform-13500000.otg_new-usb-0:1.1:1.0", KERNEL=="video*", ENV{MDEV}="sub-video-$env{MINOR}", SYMLINK+="v4l/by-id/sub-video-$env{MINOR}", RUN+="/usr/bin/auto_uvc.sh sub-video-$env{MINOR}"
# USB HUB: DM3/DP3
#ENV{ID_PATH}=="platform-13500000.otg_new-usb-0:1.3:1.0", KERNEL=="video*", ENV{MDEV}="main-video-$env{MINOR}", SYMLINK+="v4l/by-id/main-video-$env{MINOR}", RUN+="/usr/bin/auto_uvc.sh main-video-$env{MINOR}"
ENV{ID_PATH}=="platform-13500000.otg_new-usb-0:1.*", KERNEL=="video*", ENV{MDEV}="main-video-$env{MINOR}", SYMLINK+="v4l/by-id/main-video-$env{MINOR}", RUN+="/usr/bin/auto_uvc.sh main-video-$env{MINOR}"

# check for valid "index" number
TEST!="index", GOTO="persistent_v4l_end"
ATTR{index}!="?*", GOTO="persistent_v4l_end"

IMPORT{builtin}="path_id"
ENV{ID_PATH}=="?*", KERNEL=="video*|vbi*", SYMLINK+="v4l/by-path/$env{ID_PATH}-video-index$attr{index}"
ENV{ID_PATH}=="?*", KERNEL=="audio*", SYMLINK+="v4l/by-path/$env{ID_PATH}-audio-index$attr{index}"

GOTO="persistent_v4l_end"

LABEL="persistent_v4l_remove"

SUBSYSTEM!="video4linux", GOTO="persistent_v4l_end"

# USB HUB: DM1/DP1
#ENV{DEVPATH}=="/devices/platform/ahb2/13500000.otg_new/usb1/1-1/1-1.1*", KERNEL=="video*", ENV{MDEV}="sub-video-$env{MINOR}",  RUN+="/usr/bin/auto_uvc.sh sub-video-$env{MINOR}"
# USB HUB: DM3/DP3
#ENV{DEVPATH}=="/devices/platform/ahb2/13500000.otg_new/usb1/1-1/1-1.3*", KERNEL=="video*", ENV{MDEV}="main-video-$env{MINOR}", RUN+="/usr/bin/auto_uvc.sh main-video-$env{MINOR}"
ENV{DEVPATH}=="/devices/platform/ahb2/13500000.otg_new/usb1/1-1/1-1.*", KERNEL=="video*", ENV{MDEV}="main-video-$env{MINOR}", RUN+="/usr/bin/auto_uvc.sh main-video-$env{MINOR}"

LABEL="persistent_v4l_end"
